# deploy/systemd/hermes-compose.service
[Unit]
Description=Hermes AI Docker Compose stack
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/hermes-ai/deploy/compose
RemainAfterExit=yes
ExecStart=/usr/bin/docker compose --env-file hermes.env up -d
ExecStop=/usr/bin/docker compose --env-file hermes.env down

[Install]
WantedBy=multi-user.target
